# srcs/influxdb/supervisord.conf
[supervisord]
user=root
nodaemon=true

[program:influxd]
command = /usr/sbin/influxd --config /etc/influxdb.conf
autorestart = true

[program:telegraf]
command = /usr/bin/telegraf
autorestart = true
