# srcs/nginx/supervisord.conf
[supervisord]
user = root
nodaemon = true

[program:nginx]
command = /usr/sbin/nginx -g "daemon off;"
autorestart = true

[program:telegraf]
command = /usr/bin/telegraf
autorestart = true
